# ansible/roles/apachemon/templates/apachemon.conf.j2
[program:apachemon]
priority=10
directory=/usr/local/sbin
command=python apachemon.py
user=apachemon
autostart=true
autorestart=true
startsecs = 0
